# cursor-sync.service
[Unit]
Description=Cursor Code Sync Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/rust/cursor_codes/cursor_codes
ExecStart=/home/rust/cursor_codes/cursor_codes/sync.sh --daemon
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cursor-sync

[Install]
WantedBy=default.target 